# root/etc/inittab
::respawn:/usr/local/sbin/php-fpm

::respawn:/usr/sbin/nginx
